# argonfan.timer
[Unit]
Description=Timer controlling the case fan

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
