# argonfan.service
[Unit]
Description=Set case fan speed

[Service]
ExecStart=/usr/local/bin/argonfan